# monthly_proj/.config/uwsgi/uwsgi.service
[Unit]
Description=uWSGI service
After=syslog.target

[Service]
ExecStart=/usr/local/bin/uwsgi -i /app/kdt-ai3-eda-visualization/monthly_proj/.config/uwsgi/mysite.ini

Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target